# deploy/content-studio.service.example
[Unit]
Description=Content Studio Workbench
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=content-studio
Group=content-studio
WorkingDirectory=/opt/content-studio/current
EnvironmentFile=/etc/content-studio/content-studio.env
ExecStartPre=/usr/bin/npm run release:preflight
ExecStart=/usr/bin/node /opt/content-studio/current/dist-server/server.mjs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=content-studio
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/content-studio /var/backups/content-studio

[Install]
WantedBy=multi-user.target
